# supervisord.conf
[supervisord]
nodaemon=true
logfile=/app/logs/supervisord.log
pidfile=/tmp/supervisord.pid
childlogdir=/app/logs
loglevel=info

[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

; =============================================================================
; GUNICORN - Django Web Server
; =============================================================================
[program:gunicorn]
command=gunicorn config.wsgi:application --bind 0.0.0.0:8000 --workers %(ENV_GUNICORN_WORKERS)s --timeout 120 --access-logfile - --error-logfile -
directory=/app
autostart=true
autorestart=true
startsecs=5
stopwaitsecs=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=DJANGO_SETTINGS_MODULE="config.settings"

; =============================================================================
; CELERY WORKER - Background Task Processing
; =============================================================================
[program:celery_worker]
command=celery -A config worker --loglevel=info --concurrency=%(ENV_CELERY_CONCURRENCY)s -Q celery,email_processing,default
directory=/app
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=60
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=DJANGO_SETTINGS_MODULE="config.settings"

; =============================================================================
; CELERY BEAT - Scheduled Task Scheduler
; =============================================================================
[program:celery_beat]
command=celery -A config beat --loglevel=info --scheduler django_celery_beat.schedulers:DatabaseScheduler
directory=/app
autostart=true
autorestart=true
startsecs=10
stopwaitsecs=30
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=DJANGO_SETTINGS_MODULE="config.settings"

; =============================================================================
; GROUP
; =============================================================================
[group:imaradesk]
programs=gunicorn,celery_worker,celery_beat
priority=999
